# puppet/modules/jenkins/templates/jenkins.service.erb
[Unit]
Description=Jenkins Daemon
After=syslog.target

[Service]
Type=simple
Environment="JENKINS_HOME=/var/lib/jenkins"
ExecStartPre=/usr/libexec/wait4mariadb
ExecStart=/usr/bin/java -Xmx<%= @settings['jenkins_heap_limit'] %> -Dhudson.model.ParametersAction.keepUndefinedParameters=true -jar /usr/lib/jenkins/jenkins.war --webroot=/var/cache/jenkins/war --httpPort=8082 --ajp13Port=-1 --prefix=/jenkins --httpListenAddress=127.0.0.1
ExecStartPost=/usr/libexec/wait4jenkins
User=jenkins

[Install]
WantedBy=multi-user.target
