# puppet/modules/nodepool/templates/nodepool.service.erb
[Unit]
Description=Nodepool daemon
After=network.target

[Service]
Type=simple
User=jenkins
Group=jenkins
ExecStart=/usr/bin/nodepoold -p /var/run/nodepool/nodepool.pid -l /etc/nodepool/nodepool.logging.conf
Environment="STATSD_HOST=statsd.<%= @fqdn %>" "STATSD_PORT=8125"
KillSignal=SIGUSR1

[Install]
WantedBy=multi-user.target
